# deploy/systemd/palinode-embed-keepalive.timer.template
[Unit]
Description=Palinode BGE-M3 GPU keepalive (every 20 min)
After=network.target

[Timer]
OnBootSec=2min
OnUnitActiveSec=20min
Unit=palinode-embed-keepalive.service

[Install]
WantedBy=default.target
