# deploy/systemd/palinode-api.service.template
[Unit]
Description=Palinode API Server (FastAPI / uvicorn)
Documentation=https://github.com/phasespace-labs/palinode
After=network.target

[Service]
Type=simple
WorkingDirectory=${PALINODE_HOME}
Environment="PALINODE_DIR=${PALINODE_DATA_DIR}"
Environment="OLLAMA_URL=${OLLAMA_URL}"
Environment="EMBEDDING_MODEL=${EMBEDDING_MODEL}"
# Suppress the 0.0.0.0 binding warning for intentional network-exposed deployments.
# This template binds to 0.0.0.0 so the API is reachable over private VPN (or any
# other network interface). PALINODE_API_BIND_INTENT=public makes that intent
# explicit and prevents misleading startup warnings. See docs/MCP-SETUP.md.
Environment="PALINODE_API_BIND_INTENT=public"
ExecStart=${PALINODE_HOME}/venv/bin/uvicorn palinode.api.server:app --host 0.0.0.0 --port ${API_PORT}
Restart=always
RestartSec=5

# Structured JSON logs → journalctl; no custom log file needed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=palinode-api

[Install]
WantedBy=default.target
